# bin/wormhole-transfer-agent/build/supervisord.conf
[supervisord]
user=root
nodaemon=true
childlogdir=/log
logfile=/log/supervisord.log

[program:openconnect]
command=/openconnect.sh
stdout_logfile=/log/openconnect.log
redirect_stderr=true

[program:inotifywait]
command=/inotifywait.sh
stdout_logfile=/log/inotifywait.log
redirect_stderr=true